# imx/meta-fsl-arm/recipes-core/udev/udev-extraconf/mx6/9-blacklist.rules

ACTION=="add", DEVPATH=="/devices/*", ENV{MODALIAS}=="?*", RUN+="/sbin/modprobe -b $env{MODALIAS}"